# data/eos-kolibri-system-helper.socket
[Unit]
Description=Kolibri Server

[Socket]
ListenStream=18009

[Install]
WantedBy=sockets.target
