# data/eos-kolibri-system-helper.socket.in
[Unit]
Description=Kolibri Server

[Socket]
ListenStream=@KOLIBRI_PORT@
Accept=no

[Install]
WantedBy=sockets.target
